# usbmount@.service.in
[Unit]
Description=Auto Mount USB Drive on /dev/%i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@PREFIX@/usbmount add    %i
ExecStop=@PREFIX@/usbmount  remove %i
